# build/packaging/debian/DEBIAN/postrm
#!/bin/bash
set -e

case "$1" in
    remove|purge)
        echo "ðŸ—‘ï¸  AILinux TriForce - Removing..."
        systemctl stop triforce 2>/dev/null || true
        systemctl disable triforce 2>/dev/null || true
        rm -f /etc/systemd/system/triforce.service
        systemctl daemon-reload
        ;;
esac

exit 0
